# scripts/aurumbotx-hyperliquid.service
[Unit]
Description=AurumBotX Hyperliquid Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/AurumBotX
Environment="PATH=/home/ubuntu/AurumBotX/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/ubuntu/AurumBotX/venv/bin/python3 /home/ubuntu/AurumBotX/wallet_runner_hyperliquid.py
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=30s

[Install]
WantedBy=multi-user.target
